# rstudio/supervisord.ini
[supervisord]
nodaemon=true

[program:proxy]
command=python3 -u /opt/proxy.py
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/proxy.log
stderr_logfile=/var/log/proxy.err.log

[program:init]
command=/usr/lib/rstudio-server/bin/rserver --auth-none=1
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/init.log
stderr_logfile=/var/log/init.err.log
environment=USER=rstudio